# scripts/systemd/linkdownloaderbotforgroups-yt-dlp-update.timer
[Unit]
Description=Run LinkDownloaderBotForGroups yt-dlp updater nightly

[Timer]
OnCalendar=*-*-* 03:27:00
RandomizedDelaySec=30m
Persistent=true
Unit=linkdownloaderbotforgroups-yt-dlp-update.service

[Install]
WantedBy=timers.target
